# rpi-reverse-tunnel/templates/rpi-reverse-tunnel-heartbeat.service
[Unit]
Description=Send reverse tunnel heartbeat

[Service]
Type=oneshot
EnvironmentFile=/etc/default/rpi-reverse-tunnel
ExecStart=/bin/bash /opt/rpi-reverse-tunnel/scripts/tunnel_health.sh beat
